# packaging/nubix.service
[Unit]
Description=Nubix Cloud Sync Manager
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
# install.sh legt den Launcher unter ~/.local/bin ab (%h = Home des Nutzers)
ExecStart=%h/.local/bin/nubix --background
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
